# hyrule-office-api.service
[Unit]
Description=Hyrule Office Status API (Python)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/littlebot/hyrule-office
ExecStart=/usr/bin/python3 /home/littlebot/hyrule-office/api-server.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
